# scripts/carrier-pigeon-deploy.service
[Unit]
Description=Carrier Pigeon auto-deploy

[Service]
Type=oneshot
ExecStart=/opt/carrier-pigeon/scripts/deploy.sh
