# scripts/carrier-pigeon-deploy.timer
[Unit]
Description=Carrier Pigeon auto-deploy timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
